# airootfs/home/asira/.config/autostart/flameos-installer.desktop

[Desktop Entry]
Type=Application
Name=AsiraOS Installer
Exec=asiraos-installer
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
