# airootfs/home/flame/.config/autostart/flameos-installer.desktop
[Desktop Entry]
Type=Application
Name=FlameOS Installer
Exec=flameos-installer
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true

